# crontab/act-esa-incidents-geojson
# update every 5 minutes
*/5 *  *   *   *  root  act-esa-incidents-geojson /srv/www/act-esa-incidents.geojson
